# deploy/systemd/tmjconnect-api-test.service
[Unit]
Description=TMJConnect API (test)
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=tmjconnect-test
Group=tmjconnect-test
WorkingDirectory=/opt/tmjconnect-test/apps/api
EnvironmentFile=/etc/tmjconnect/api-test.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tmjconnect-api-test

# ─── Hardening ───────────────────────────────────────────────────────────────
# Same surface as prod — test is a behavioral mirror, not a relaxed sandbox.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/var/lib/tmjconnect-test
AmbientCapabilities=
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
